# deploy/divi-rpc-proxy-mainnet.service
[Unit]
Description=DIVI RPC Filtering Proxy (Mainnet)
After=network.target privatedivi-mainnet-node1.service
Wants=privatedivi-mainnet-node1.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/divi-rpc-proxy \
    --listen-addr 127.0.0.1 \
    --listen-port 18081 \
    --backend-url http://127.0.0.1:52491 \
    --rpc-user privatedivi \
    --rpc-password ${RPC_PASSWORD_MAINNET} \
    --network mainnet

EnvironmentFile=/etc/divi-rpc-proxy/mainnet.env
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
